# deploy/quadlet/hermes-default.container
[Unit]
Description=Hermes Agent — default bot (sandbox / scratch)
# Pattern: one container per bot (Telegram persona). N bots = N Quadlets sharing
# the same image, each with its own ~/.hermes/bots/<name>/ data dir. Container-level
# isolation replaces hermes' built-in `profiles` system — simpler, granular restart,
# isolated logs, OOM blast-radius capped per bot.
# Outbound-only (LLM API + chat platform webhooks/polling); network-online.target
# is enough. No dependency on other hermes-* services.
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Container]
# Official upstream image on Docker Hub. AutoUpdate=registry → podman-auto-update.timer
# checks the digest every 5 min and restarts on change (see lyra
# docs/ops/container-publishing.md § Auto-Update Flow).
# To pin a known-good version, replace `:latest` with `:vX.Y.Z` and remove the Label.
Image=docker.io/nousresearch/hermes-agent:latest
# Label=io.containers.autoupdate=registry  # DISABLED 2026-06-04: docker.io anon
# rate-limit (toomanyrequests) poisoned the whole `podman auto-update` batch on M1,
# exit-125 every run -> factory-* never auto-pulled. To restore: re-add this Label +
# `podman login docker.io` (or pin Image=:vX.Y.Z and bump manually).
ContainerName=hermes-default
Network=roxabi.network

# Hardening — defense in depth (matches lyra invariants where compatible).
NoNewPrivileges=true
# DEVIATION from lyra `DropCapability=all`: upstream s6-overlay requires
# CAP_SETUID to drop from root to the hermes user at startup. Without it,
# `s6-overlay-suexec` fails with `unable to set supplementary group list` or
# `/run belongs to uid 0 instead of 10000` (depends on digest). Verified
# 2026-05-28 on digest 67011cc0…: `DropCapability=all` + `User=10000:10000`
# produces a crashloop. We keep `NoNewPrivileges=true` as a backstop.
# DropCapability=all

# DEVIATION from lyra `ReadOnly=true`: hermes ships a writable /opt/hermes/.venv
# by design — `tools/lazy_deps.py` installs platform-specific adapter packages
# (discord.py, etc.) at first use, and the Dockerfile explicitly chowns .venv to
# the hermes user with this contract documented. Read-only root FS breaks runtime
# `uv pip install`. The `--extra all` build covers most paths; if you confirm
# none of your adapters trigger lazy_deps in your config, you can enable
# ReadOnly=true and accept the residual failure mode.
# ReadOnly=true

# UserNS=keep-id maps container UID 10000 (the baked-in `hermes` user) to host
# UID 1000 (mickael), which owns the bot data dir. We force User=0:0 so the
# container starts as root — s6-overlay needs root to chown /run and related
# s6 dirs, then it drops to hermes via `su` internally. The Quadlet generator
# auto-injects --user 10000:10000 when UserNS=keep-id:uid=10000,gid=10000 is
# used without an explicit User=, which breaks the s6 bootstrap on recent
# digests (2026-05-28). Hermes still runs as uid 10000 inside the container;
# only the brief s6-init phase is root.
UserNS=keep-id:uid=10000,gid=10000
User=0:0

# Bind mount: ~/.hermes/bots/default ↔ /opt/data. Holds this bot's config.yaml,
# .env (Telegram token + LLM keys), sessions/, skills/, memories/, plans/, cron/,
# logs/, hooks/, home/, workspace/. Fully isolated from other hermes-* bots.
# Image is stateless — upgrade by `podman pull` (or auto-update); state survives.
Volume=hermes-default.volume:/opt/data:z

# Protect .env and config.yaml against upstream image seeding (upstream
# stage2-hook.sh copies .env.example / cli-config.yaml.example over existing
# files on some digests — see 2026-05-28 incident). The :ro over-mount masks
# the files inside the bind so the init script cannot clobber them.
# Hermes itself only reads these files; it does not need write access.
Volume=/home/mickael/.hermes/bots/default/.env:/opt/data/.env:ro
Volume=/home/mickael/.hermes/bots/default/config.yaml:/opt/data/config.yaml:ro

# No PublishPort in pure gateway mode — Telegram polling, Discord WS, Slack
# webhooks, LLM API calls are all outbound. Uncomment ONE of the blocks below
# only if you need host access; see security notes in docker-compose.yml.
#
# OpenAI-compat API server (requires API_SERVER_KEY ≥ 8 chars):
# Environment=API_SERVER_ENABLED=true
# Environment=API_SERVER_HOST=0.0.0.0
# Environment=API_SERVER_KEY=
# PublishPort=127.0.0.1:8642:8642
#
# Web dashboard (in-container side-process, not auto-restarted if it crashes):
# Environment=HERMES_DASHBOARD=1
# PublishPort=127.0.0.1:9119:9119

Exec=gateway run

[Service]
Restart=on-failure
RestartSec=5
# Cap at 2 cores — voice transcription + LLM streaming can spin up workers.
CPUQuota=200%

[Install]
WantedBy=default.target
